# templates/default/systemd.service.erb
[Unit]
Description=IAP auth service
Documentation=https://github.com/gojekfarm/iap_auth
ConditionPathExists=<%= @iap_auth_dir %>/iap_auth
After=network.target

[Service]
WorkingDirectory=<%= @iap_auth_dir %>
ExecStart=<%= @iap_auth_dir %>/iap_auth server
Restart=always
RestartSec=2
StartLimitInterval=50
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
